# scripts/loop-drop.service
[Unit]
Description=Loop Drop
After=multi-user.target

[Service]
Type=idle
ExecStart=/home/pi/loop-drop

[Install]
WantedBy=multi-user.target