# scripts/loop-drop.service
[Unit]
Description=Loop Drop
After=local-fs.target

[Service]
Type=idle
ExecStart=/home/pi/loop-drop

[Install]
WantedBy=local-fs.target